# contrail/systemctl/contrail@api.service
[Unit]
Description = Contrail contrail-api

[Service]
ExecReload = /bin/kill -HUP $MAINPID
TimeoutStopSec = 300
KillMode = process
ExecStart = /home/cloud/.nix-profile/bin/contrail-api --conf_file /etc/contrail/contrail-api.conf
User = cloud

[Install]
WantedBy = multi-user.target
